# data/grafana/periodic/daily/grafana-backup-daily.sh
#!/bin/sh

sh /etc/scripts/backup-file.sh daily 5 /var/lib/grafana/grafana.db

