# data/grafana/periodic/hourly/grafana-backup-hourly.sh
#!/bin/sh

sh /etc/scripts/backup-file.sh hourly 8 /var/lib/grafana/grafana.db
